# device-configs/etc/sv/nocturned/run
#!/bin/sh

exec 2>&1
sv check dbus > /dev/null || exit 1
sv check bluetoothd > /dev/null || exit 1  
sv check bluetooth_adapter > /dev/null || exit 1

# Wait for filesystems to be fully mounted and writable
# This is critical for car power scenarios where boot is slower
echo "Waiting for filesystems..."
for i in 1 2 3 4 5 6 7 8 9 10; do
    if touch /var/.test 2>/dev/null; then
        rm -f /var/.test
        echo "Filesystem ready after $i seconds"
        break
    fi
    sleep 1
done

# Give system time to stabilize after filesystem mount
# This helps with power-constrained environments
sleep 3

# Give Bluetooth hardware extra time to fully initialize after btattach
sleep 5

# Wait for Bluetooth adapter to be fully ready
for i in 1 2 3 4 5; do
    if bluetoothctl show > /dev/null 2>&1; then
        break
    fi
    sleep 2
done

# Create log directory if it doesn't exist
mkdir -p /var/nocturne/logs

# Start nocturned with output for debugging
# Redirect to log file and let runsv capture it
exec /usr/local/bin/nocturned 2>&1
