# device-configs/etc/sv/nocturned/run
#!/bin/sh

exec 2>&1
sv check dbus > /dev/null || exit 1
sv check bluetoothd > /dev/null || exit 1  
sv check bluetooth_adapter > /dev/null || exit 1

# Quick filesystem check (reduced from 10s to 3s max)
for i in 1 2 3; do
    if touch /var/.test 2>/dev/null; then
        rm -f /var/.test
        break
    fi
    sleep 1
done

# Minimal stabilization delay
sleep 1

# Wait for Bluetooth adapter (reduced from 10s to 3s max)
for i in 1 2 3; do
    if bluetoothctl show > /dev/null 2>&1; then
        break
    fi
    sleep 1
done

# Create log directory if it doesn't exist
mkdir -p /var/nocturne/logs

# Start nocturned with output for debugging
# Redirect to log file and let runsv capture it
exec /usr/local/bin/nocturned 2>&1
